# www.conf
[program:www]
command=python3 /home/curator2/www/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/www/test.err.log
stdout_logfile=/home/curator2/www/test.out.log